# tools/tomee_srv/supervisord.conf
[supervisord]
nodaemon=true

[program:tomee]
command=/opt/tomee/bin/catalina.sh run

