# chat/proxy/supervisord.conf
[supervisord]
nodaemon=true

[program:proxy]
command=/usr/sbin/nginx

